# deploy/dnc-fleetd.service
# divide-and-conquer fleet daemon — inventory, discovery, IaC plays, task ledger.
#   sudo cp ~/dnc/systemd/dnc-fleetd.service /etc/systemd/system/   # rendered by bootstrap.sh (templates use __DNC_*__)
#   sudo systemctl daemon-reload && sudo systemctl enable --now dnc-fleetd
# Adjust User/paths if your layout differs from the bootstrap defaults.

[Unit]
Description=divide-and-conquer fleet daemon (fleetd)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__DNC_USER__
WorkingDirectory=__DNC_HOME__
EnvironmentFile=__DNC_HOME__/.env
ExecStart=__DNC_HOME__/.venv/bin/fleetd serve --host 0.0.0.0 --port 7431
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
